# php-fpm/supervisor/supervisord.conf
[supervisord]
[program:worker]
command=php /var/www/html/web/yii cdn-crop-queue/listen
process_name=%(program_name)s_%(process_num)02d
numprocs=proc_count
stdout_logfile=/var/log/worker.log
autostart=true
autorestart=true
user=www-data
stopsignal=KILL